# skel/etc/supervisor/conf.d/900-logstash-forwarder.conf
[program:logstash-forwarder]
command=%(ENV_SERVICE_LOGSTASH_FORWARDER_CMD)s
user=logstash-forwarder
autostart=true
autorestart=unexpected
startsecs=5
priority=900
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
