# skel/etc/supervisor/conf.d/999-redpill.conf
[program:redpill]
command=%(ENV_SERVICE_REDPILL_CMD)s
autostart=true
priority=999
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
